# mkosi/mkosi.extra/usr/lib/systemd/system/kick-scan.service
[Unit]
After=multi-user.target

[Service]
Type=oneshot
ExecStart=iw wlp6s0 scan

[Install]
WantedBy=multi-user.target
