# roles/ansible_pull/templates/ansible-pull-repeat.timer.j2
[Unit]
Description=Run ansible setup after repo change

[Timer]
{{ ansible_pull_schedule_repeat }}

[Install]
WantedBy=timers.target
