# roles/ansible_pull/templates/ansible-pull-init.timer.j2
[Unit]
Description=Run full ansible setup

[Timer]
{{ ansible_pull_schedule_init }}

[Install]
WantedBy=timers.target
