# roles/ansible-pull/templates/ansible-pull-init.timer.j2
[Unit]
Description=Run full ansible setup

[Timer]
{{ schedule_init }}

[Install]
WantedBy=timers.target
