# roles/ansible-pull/templates/ansible-pull-repeat.timer.j2
[Unit]
Description=Run ansible setup after repo change

[Timer]
{{ schedule_repeat }}

[Install]
WantedBy=timers.target
